# pre_order_bot/config/bot.conf
[program:<наименование процесса>]
command=<путь до пайтон> <путь до исполняемого файла>
directory=<корневая директория проекта>
user=<имя пользователя>
autorestart=true
redirect_stderr=true